# services-s6/.local/share/s6/sv/rill/run
#!/bin/sh
cd "$HOME" || exit 1
exec rill
