# service-turnstile-runit/.config/service/wlsunset/run
#!/bin/sh
exec 2>&1
exec wlsunset -l 40 -L -8 -t 3500 -T 5000 -g 0.8
