# services-systemd/.config/systemd/user/graphical-session.target
[Unit]
Description=User Graphical Session
Requires=river.service
After=river.service
