# services-s6/.local/share/s6/sv/dbus-session-log/run
#!/bin/sh
d="${XDG_STATE_HOME:-$HOME/.local/state}/s6-logs/dbus-session"
mkdir -p "$d"
exec s6-log -d3 t "$d"
