# service-turnstile-runit/.config/service/gammastep/run
#!/bin/sh
exec 2>&1
exec gammastep -l 40:-8 -t 6500:4500 -g 0.8
