# services-runit-turnstile/.config/service/footserver/run
#!/bin/sh
exec 2>&1
sv check ../wayland-ready || exit 1
[ -r "$HOME/.profile" ] && . "$HOME/.profile"
exec chpst -e "$TURNSTILE_ENV_DIR" foot --server
